# mrstone-website.service
[Unit]
Description=Mr Stone Personal Website
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/docusaurus/chenxiaoshivivid
ExecStartPre=/usr/bin/npm install
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/npm run serve
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=mrstone-website
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 